# roles/java/files/etc/profile.d/java.sh
export JAVA_HOME="/usr/lib/jvm/jdk8u282-b08"
export PATH="$JAVA_HOME/bin:$PATH"
